# tieba-sign.service
[Unit]
After = network-online.target

[Service]
ExecStart = tieba-sign
Restart = no
Type = oneshot
